# _developer-directory-dashboard_/clearlinux/base/usr/lib/systemd/system/clr-elf-replace-trigger.service
[Unit]
Description=Replace elf files with optimzed versions
BindsTo=update-triggers.target
ConditionPathExists=/usr/libexec/clr-elf-replace

[Service]
Type=oneshot
ExecStart=/usr/libexec/clr-elf-replace
