# _developer-directory-dashboard_/clearlinux/base/usr/lib/systemd/system/dynamic-trust-store.service
[Unit]
Description=Install prebuilt store if trust management is not installed.
ConditionFileIsExecutable=/usr/bin/clrtrust

[Service]
Type=oneshot
ExecStart=/usr/bin/clrtrust generate

[Install]
WantedBy=update-triggers.target
