# _developer-directory-dashboard_/clearlinux/base/usr/lib/systemd/system/firstboot-triggers.service
[Unit]
Description=Update triggers on first boot
ConditionFirstBoot=|yes
ConditionPathExists=|!/var/cache/ldconfig/ld.so.cache

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/systemctl start --no-block update-triggers.target
