# _developer-directory-dashboard_/clearlinux/base/usr/lib/systemd/system/networkd-fallback.service
[Unit]
Description=Attempt to start systemd-networkd if no one manages ethernet if

[Service]
Type=oneshot
ExecStart=/usr/libexec/clr-fallback-to-systemd-networkd
