# _developer-directory-dashboard_/clearlinux/base/usr/lib/systemd/system/networkd-fallback.timer
[Unit]
Description=Attempt to start systemd-networkd if needed 1 min after boot

[Timer]
OnStartupSec=60s

[Install]
WantedBy=timers.target
