# _developer-directory-dashboard_/clearlinux/base/usr/lib/systemd/system/pacdiscovery.path
[Unit]
Description=Proxy AutoConfig discovery service (lease directory watch)

[Path]
DirectoryNotEmpty=/run/systemd/netif/leases
DirectoryNotEmpty=/run/NetworkManager/system-connections

[Install]
WantedBy=multi-user.target
