# _developer-directory-dashboard_/clearlinux/base/usr/lib/systemd/system/pacrunner.service
[Unit]
Description=Proxy AutoConfig runner service

[Service]
ExecStart=/usr/sbin/pacrunner -n
User=pacrunner

[Install]
WantedBy=multi-user.target
