# _developer-directory-dashboard_/clearlinux/base/usr/lib/systemd/system/swupd-cleanup.service
[Unit]
Description=Clean up swupd state directory

[Service]
Type=oneshot
ExecStart=/usr/bin/swupd clean
