# _developer-directory-dashboard_/clearlinux/base/usr/lib/systemd/system/sysusers-trigger.service
[Unit]
Description=Create sysusers
Documentation=man:systemd-sysusers(8)
BindsTo=update-triggers.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemd-sysusers
